# plugins/kalabox-syncthing/dockerfiles/syncthing/syncthing-supervisor.conf
[program:syncthing]
command = /usr/local/bin/syncthing -no-browser -home /syncthing/
directory = /sync/
autorestart = True
user = root
