# plugins/kalabox-syncthing/dockerfiles/syncthing/syncthing-supervisor.conf
[program:syncthing]
command = /usr/local/bin/syncthing -no-browser -home /etc/syncthing/
directory = /root/
autorestart = True
user = root
