# dockerfiles/kalabox/syncthing/syncthing-supervisor.conf
[program:syncthing]
command = /usr/local/bin/syncthing -no-browser -home /etc/syncthing/
directory = /sync/
autorestart = True
user = root
