# dockerfiles/syncthing/syncthing-supervisor.conf
[program:syncthing]
command = /usr/local/bin/syncthing -no-browser -home /data/
directory = /sync/
autorestart = True
user = root
